# root/etc/uci-defaults/80_aurora
#!/bin/sh

readonly CONF_FILE="/etc/config/aurora"
readonly TEMPLATE_FILE="/usr/share/aurora/classic.template"
readonly LOG_TAG="aurora"

log() {
    logger -t "$LOG_TAG" "$1"
}

main() {
    log "Initializing Aurora theme configuration"

    if [ ! -f "$TEMPLATE_FILE" ]; then
        log "Template missing: $TEMPLATE_FILE"
        exit 0
    fi

    if [ ! -f "$CONF_FILE" ]; then
        if cp "$TEMPLATE_FILE" "$CONF_FILE"; then
            chmod 644 "$CONF_FILE"
            log "First-time initialization: Config created from template"
        else
            log "Failed to create config from template"
        fi
    else
        log "Configuration already exists, skipping initialization"
    fi
}

main

exit 0
